# documentation/docker_images/base_image/configs/etc/cron.d/genular
# /etc/cron.d/genular: crontab fragment for genular
# Look for new analysis and try to process them

* * * * * root /usr/bin/flock -n /tmp/simon_cron.lock /usr/bin/Rscript /var/www/genular/simon-backend/cron/main.R cron >>/var/log/simon-cron.log 2>&1
